# packages/cleanup-journald/cleanup-journald.service
[Unit]
Description=Check for Low Battery

[Service]
Type=oneshot
ExecStart=/usr/bin/journalctl --vacuum-time=1m
