# packages/cleanup-journald/cleanup-journald.service
[Unit]
Description=Cleanup old journald logs (older than 1 month)

[Service]
Type=oneshot
ExecStart=/usr/bin/journalctl --vacuum-time=1m
