# packages/cleanup-journald/cleanup-journald.timer
[Unit]
Description=Cleanup old journald logs daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
