# packages/cleanup-journald/cleanup-journald.service
[Unit]
Description=Cleanup journald logs daily

[Service]
Type=oneshot
ExecStart=/usr/bin/journalctl --vacuum-time=1m
