# packages/reflector-service-hooks/reflector.service
[Unit]
Description=Pacman mirrorlist update
Documentation=https://wiki.archlinux.org/index.php/Reflector
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/reflector --age 24 --country US --latest 300 --number 50 --sort rate --save /etc/pacman.d/mirrorlist -p https -p http -p rsync -p ftp
